# nullifier-ingest/docs/nullifier-query-server.service
[Unit]
Description=Nullifier PIR server (nf-server)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nullifier-ingest
ExecStart=/opt/nullifier-ingest/nf-server serve --pir-data-dir /opt/nullifier-ingest/pir-data --port 3000
Restart=on-failure

[Install]
WantedBy=multi-user.target
