# sdk/docs/zallyd-val2.service
[Unit]
Description=Zally SDK chain - Validator 2 (zallyd-val2)
After=network.target zallyd-val1.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/zally-chain
ExecStart=/opt/zally-chain/zallyd start --home /opt/zally-chain/.zallyd-val2
Restart=on-failure

[Install]
WantedBy=multi-user.target
